# systemd/adx-dashboard.service
[Unit]
Description=ADX Strategy v2.0 - Web Dashboard
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/dev/trading/adx_strategy_v2
ExecStart=/var/www/dev/trading/adx_strategy_v2/venv/bin/python3 /var/www/dev/trading/adx_strategy_v2/dashboard_web.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=adx-dashboard

# Environment
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target
